# supervisor/mysql.conf
[program:mysql]
directory=/opt/mysql
command=/usr/bin/pidproxy /var/run/mysqld/mysqld.pid /opt/mysql/run_mysql.sh
