# supervisor/mysql.conf
[program:mysql]
directory=/opt/mysql
command=/opt/mysql/run_mysql.sh
